# .server/supervisor/manifest-worker.conf
[program:manifest-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sistem-manifest-backend/artisan queue:work --queue=scans,compliance,default --sleep=3 --tries=3 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/sistem-manifest-backend/storage/logs/worker.log
stopwaitsecs=120
